# websocket.conf
[program:novatopro]
command=/usr/bin/php /var/www/novato.pro/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=ubuntu